# unused_archive/com.auto2fa.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.auto2fa.daemon</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/shgao/logs/auto2fa_dev/auto2fa_daemon.py</string>
        <string>--daemon</string>
        <string>--config</string>
        <string>/Users/shgao/logs/auto2fa_dev/daemon_config.json</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/shgao/logs/auto2fa_dev</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/auto2fa_daemon.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/auto2fa_daemon.log</string>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ExitTimeOut</key>
    <integer>30</integer>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>PYTHONPATH</key>
        <string>/Users/shgao/logs/auto2fa_dev</string>
    </dict>
</dict>
</plist>
