# LaunchAgents/com.auto2fa.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.auto2fa.daemon</string>

    <!-- Path to the auto2fa-daemon entry point that `pip install -e .` creates.
         If your `pip install` used a non-standard prefix (homebrew, asdf,
         pyenv, conda), update this path. Find it with `which auto2fa-daemon`. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/auto2fa-daemon</string>
    </array>

    <!-- Auto-start on login -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- The daemon needs SSH_CONFIG_PATH to find passwords.json + tunnels.json.
         Adjust this to your ~/.ssh equivalent. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>SSH_CONFIG_PATH</key>
        <string>/Users/YOUR_USERNAME/.ssh</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Log to a known location for debugging -->
    <key>StandardOutPath</key>
    <string>/tmp/auto2fa_daemon.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/auto2fa_daemon.err</string>

    <!-- Reasonable resource limits -->
    <key>ProcessType</key>
    <string>Interactive</string>

    <!-- Don't restart faster than every 10s if it keeps crashing -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
